# config/systemd/fireworks-worker.target
[Unit]
StopWhenUnneeded=true
Wants=fireworks-worker-1.service fireworks-worker-2.service
