# config/systemd/fireworks-web.target
[Unit]
StopWhenUnneeded=true
Wants=fireworks-web-1.service
