# config/systemd/fireworks-scheduler.target
[Unit]
StopWhenUnneeded=true
Wants=fireworks-scheduler-1.service
