# config/systemd/fireworks.target
[Unit]
StopWhenUnneeded=true
Wants=fireworks-web.target fireworks-scheduler.target fireworks-worker.target

[Install]
WantedBy=multi-user.target
